# preset/overlay/etc/udev/rules.d/70-hdmi.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/bin/su - cqutprint -c /usr/local/bin/hdmi-toggle"
